# install/etc/s6/services/10-unbound/run
#!/usr/bin/with-contenv bash

source /assets/functions/*
PROCESS_NAME="unbound"
check_initialized
liftoff
print_info "Starting Unbound"
exec unbound -d -v
